# assignment_04/42-usb_keyboard_module.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/sbin/insmod /root/test/assignment_04/usb_keyboard_module.ko"
